# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery notification service

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/scripts/battery-notify.sh
